# helper/supervisord.conf
[program:rqworkers]
command=/home/pranavgupta/envs/aws/bin/rq worker high default low queue1 queue2 queue3 queue4 queue5
numprocs=10
process_name=%(process_num)s
autostart=true

[supervisord]