# supervisord/supervisord.conf
[program:rqworkers]
command=/home/pranavgupta/envs/aws/bin/rq worker high default low
numprocs=10
process_name=%(process_num)s
autostart=true

[supervisord]